# .config/systemd/user/lowbat.timer
[Unit]
Description=run lowbat every minute

[Timer]
OnCalendar=minutely
Persistent=true
