# .config/systemd/user/lowbat.service
[Unit]
Description=low battery warnings
Wants=lowbat.timer

[Service]
ExecStart=/bin/bash %h/lowbat.sh

[Install]
WantedBy=default.target
