# .config/systemd/user/grive-timer@grive.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/lib/grive/grive-sync.sh sync "%i"