# .config/systemd/user/grive-changes@grive.service.d/override.conf
[Service]
ExecStart=/usr/lib/grive/grive-sync.sh listen "%i"